# etc/systemd/system/dn42-roa.timer
[Unit]
Description=Periodically update DN42 roa

[Timer]
OnBootSec=15s
OnActiveSec=1h

[Install]
WantedBy=dn42.target

