# etc/systemd/system/dn42.target
[Unit]
Description=DN42 target
Requires=network.target

[Install]
WantedBy=multi-user.target
